# supervisor/compy.conf
[program:compy]
command=/opt/compy/docker.sh
