# supervisor/http.conf
[program:http]
command=sed -i -e "s|SERVER_CHANGE_ME|$HOST|g" /opt/compy/proxy.pac && http-server /opt/compy
