# supervisor/http.conf
[program:http]
command="/usr/local/bin/http-server ."
